# popopop.service
[Unit]
Description=POPoPOP.Top - POP Text Editor Web App
Documentation=https://github.com/shinyzhu/popopop.top
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/popopop.top
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=popopop

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/popopop.top/public/fonts
PrivateTmp=true

[Install]
WantedBy=multi-user.target
